# scripts/com.skillforge.prcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.skillforge.prcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ahnjundaram_g/dev/tools/skill-forge/scripts/check-pr-status.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>21600</integer><!-- 6시간마다 (초 단위) -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/skillforge-prcheck.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/skillforge-prcheck.log</string>
</dict>
</plist>
